# archive_compression/crontab
0 2 * * * python /archive_compression.py >> /logfile
