# archive_compression/crontab
0 2 * * * python3 /archive_compression.py >> /var/log/archivecompression.log
